# systemd/pm2.service
[Unit]
Description=Restart PM2 processes

[Service]
Type=oneshot
ExecStart=/bin/bash /home/admin/projectlab/app/scripts/startup.sh