# systemd/wos-sync.timer
[Unit]
Description=Timer for running syncing with Wizeline OS once a week

[Timer]
OnUnitActiveSec=7d
OnBootSec=10s

[Install]
WantedBy=timers.target